# bot.conf
[program:telegrambot]
directory=/home/ubuntu/Downloader
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/Downloader/logfile_err.log
stdout_logfile=/home/ubuntu/Downloader/logfile.log